# contrib/init/apixcashd.conf
description "Apixcash Core Daemon"

start on runlevel [2345]
stop on starting rc RUNLEVEL=[016]

env APIXCASHD_BIN="/usr/bin/apixcashd"
env APIXCASHD_USER="apixcash"
env APIXCASHD_GROUP="apixcash"
env APIXCASHD_PIDDIR="/var/run/apixcashd"
# upstart can't handle variables constructed with other variables
env APIXCASHD_PIDFILE="/var/run/apixcashd/apixcashd.pid"
env APIXCASHD_CONFIGFILE="/etc/apixcash/apixcash.conf"
env APIXCASHD_DATADIR="/var/lib/apixcashd"

expect fork

respawn
respawn limit 5 120
kill timeout 60

pre-start script
    # this will catch non-existent config files
    # apixcashd will check and exit with this very warning, but it can do so
    # long after forking, leaving upstart to think everything started fine.
    # since this is a commonly encountered case on install, just check and
    # warn here.
    if ! grep -qs '^rpcpassword=' "$APIXCASHD_CONFIGFILE" ; then
        echo "ERROR: You must set a secure rpcpassword to run apixcashd."
        echo "The setting must appear in $APIXCASHD_CONFIGFILE"
        echo
        echo "This password is security critical to securing wallets "
        echo "and must not be the same as the rpcuser setting."
        echo "You can generate a suitable random password using the following "
        echo "command from the shell:"
        echo
        echo "bash -c 'tr -dc a-zA-Z0-9 < /dev/urandom | head -c32 && echo'"
        echo
        echo "It is recommended that you also set alertnotify so you are "
        echo "notified of problems:"
        echo
        echo "ie: alertnotify=echo %%s | mail -s \"Apixcash Alert\"" \
            "admin@foo.com"
        echo
        exit 1
    fi

    mkdir -p "$APIXCASHD_PIDDIR"
    chmod 0755 "$APIXCASHD_PIDDIR"
    chown $APIXCASHD_USER:$APIXCASHD_GROUP "$APIXCASHD_PIDDIR"
    chown $APIXCASHD_USER:$APIXCASHD_GROUP "$APIXCASHD_CONFIGFILE"
    chmod 0660 "$APIXCASHD_CONFIGFILE"
end script

exec start-stop-daemon \
    --start \
    --pidfile "$APIXCASHD_PIDFILE" \
    --chuid $APIXCASHD_USER:$APIXCASHD_GROUP \
    --exec "$APIXCASHD_BIN" \
    -- \
    -pid="$APIXCASHD_PIDFILE" \
    -conf="$APIXCASHD_CONFIGFILE" \
    -datadir="$APIXCASHD_DATADIR" \
    -disablewallet \
    -daemon

